# directories/system/etc/profile.d/bash_completion.sh
if [ -f /usr/share/bash-completion/bash_completion ]
 then . /usr/share/bash-completion/bash_completion
fi
