# directories/system/etc/profile.d/dircolors.sh
if [ -f "/etc/dircolors" ]; then eval $(dircolors -b /etc/dircolors); fi

alias ls='ls --color=auto'
alias grep='grep --color=auto'
