# ainit-sysv-deprecated/ainit-sysv-mk-ii/etc/inittab
id:3:initdefault:

si::sysinit:/etc/apathy-init/rc S

l0:0:wait:/etc/apathy-init/rc 0
l1:S1:wait:/etc/apathy-init/rc 1
l2:2:wait:/etc/apathy-init/rc 2
l3:3:wait:/etc/apathy-init/rc 3
l4:4:wait:/etc/apathy-init/rc 4
l5:5:wait:/etc/apathy-init/rc 5
l6:6:wait:/etc/apathy-init/rc 6

ca:12345:ctrlaltdel:/sbin/shutdown -t1 -a -r now

su:S016:once:/sbin/sulogin

1:2345:respawn:/sbin/agetty --noclear tty1 9600
2:2345:respawn:/sbin/agetty tty2 9600
3:2345:respawn:/sbin/agetty tty3 9600
